# windows/WSL/etc/profile.d/docker.sh
#!/bin/sh
# This is not needed when running docker in native WSL2 mode.
export DOCKER_HOST=tcp://localhost:2375
