# windows/WSL/etc/profile.d/x11-wsl1.sh
#!/bin/sh
export DISPLAY=localhost:0.0
export LIBGL_ALWAYS_INDIRECT=1
export NO_AT_BRIDGE=1
